# deb/postinst
#!/bin/bash

systemctl enable docker-registry.service
systemctl restart docker-registry.service
