# deb/prerm
#!/bin/bash

systemctl stop docker-registry.service
systemctl disable docker-registry.service
